# copyright-manager.service
# ============================================================
# systemd 服务配置文件
# 说明：用于开机自启动和管理Gunicorn进程
# 安装方法：
#   1. 复制此文件到 /etc/systemd/system/
#   2. 修改 WorkingDirectory 为实际项目路径
#   3. 修改 User 和 Group 为实际用户
#   4. 运行：sudo systemctl daemon-reload
#   5. 运行：sudo systemctl enable copyright-manager
#   6. 运行：sudo systemctl start copyright-manager
#
# 常用命令：
#   启动：sudo systemctl start copyright-manager
#   停止：sudo systemctl stop copyright-manager
#   重启：sudo systemctl restart copyright-manager
#   状态：sudo systemctl status copyright-manager
#   日志：sudo journalctl -u copyright-manager -f
# ============================================================

[Unit]
Description=SMPH外版图书管理系统 - Gunicorn服务
After=network.target mysql.service
Wants=mysql.service

[Service]
# 用户和组（请修改为实际用户）
User=www-data
Group=www-data

# 工作目录（请修改为实际路径）
WorkingDirectory=/path/to/copyright-manager

# 环境变量文件（可选，用于存储敏感配置）
EnvironmentFile=/path/to/copyright-manager/.env

# 或者直接设置环境变量
# Environment="DB_HOST=localhost"
# Environment="DB_PORT=3306"
# Environment="DB_USER=copyright_user"
# Environment="DB_PASSWORD=your_password"
# Environment="DB_NAME=copyright_manager"
# Environment="ADMIN_TOKEN=your_secure_token"

# 启动命令
ExecStart=/usr/bin/gunicorn -c gunicorn.conf.py app:app

# 重启策略
Restart=always
RestartSec=5

# 标准输出和错误输出
StandardOutput=journal
StandardError=journal

# 进程类型
Type=notify

# 优雅关闭超时
TimeoutStopSec=30

# 安全设置
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
